# systemd/oats.target
[Unit]
Description=Starts all relevant scripts and services for OATS
Requires=napalm-logs.service oatsclient.service
